# conf/magic3.ini
[program:programme]
numprocs = 14
process_name = %(program_name)s_%(process_num)02d
command = /bin/cat
weight = 1.0

