# conf/magic2.ini
[program:theprogramname]
autoreload = True
autorestart = True
weight = 10.0
process_name = %(program_name)s
command = /bin/cat
hs_minprocs = 2

