# conf/magic1.ini
[program:theprogramname]
autoreload = True
autorestart = True
numprocs = 2
weight = 10.0
process_name = %(program_name)s_%(process_num)02d
command = /bin/cat

[program:programme]
process_name = %(program_name)s_%(process_num)02d
numprocs = 2
command = /bin/cat
weight = 1.0

