# conf/magic1.ini
[program:theprogramname]
autoreload = True
autorestart = True
weight = 10.0
process_name = %(program_name)s
command = /bin/cat
hs_minprocs = 2

[program:programme]
numprocs = 3
process_name = %(program_name)s_%(process_num)02d
command = /bin/cat
weight = 1.0

