# conf/magic1.ini
[program:printotherstuff]
numprocs = 2
process_name = %(program_name)s_%(process_num)02d
command = python /home/farthvader/Projects/superwiser/printotherstuff.py
weight = 10.0

[program:printstuff]
numprocs = 3
process_name = %(program_name)s_%(process_num)02d
command = python /home/farthvader/Projects/superwiser/printstuff.py
weight = 5.0

